# etc/examples/systemd/onetimesecret-worker.service
# Onetime Secret - Background Worker Service
#
# Required only for full authentication mode (PostgreSQL + RabbitMQ).
# Processes background jobs (email delivery, account operations, etc.)
#
# Install:
#   sudo cp etc/examples/systemd/onetimesecret-*.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now onetimesecret-worker

[Unit]
Description=Onetime Secret Background Worker
Documentation=https://docs.onetimesecret.com/en/self-hosting/
After=network.target redis-server.service valkey.service onetimesecret-web.service

[Service]
Type=simple
User=onetime
Group=onetime
WorkingDirectory=/var/lib/onetimesecret

ExecStart=/bin/bash -lc 'source .env.sh && exec bin/ots worker'

# Allow time for in-flight jobs to drain
TimeoutStopSec=60

Restart=on-failure
RestartSec=5

EnvironmentFile=/etc/default/onetimesecret

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/onetimesecret/tmp /var/lib/onetimesecret/log
PrivateTmp=true

[Install]
WantedBy=multi-user.target
