# etc/examples/systemd/onetimesecret-web.service
# Onetime Secret - Web (Puma) Service
#
# Install:
#   sudo cp etc/examples/systemd/onetimesecret-*.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now onetimesecret-web
#
# Assumes:
#   - App installed to /var/lib/onetimesecret (adjust WorkingDirectory)
#   - Dedicated 'onetime' user owns the app directory
#   - Ruby available in onetime user's PATH (via ruby-install, rbenv, etc.)
#   - .env.sh sources .env and exports variables
#   - Redis/Valkey already running (see After/Requires)

[Unit]
Description=Onetime Secret Web Server (Puma)
Documentation=https://docs.onetimesecret.com/en/self-hosting/
After=network.target redis-server.service valkey.service
Requires=network.target

[Service]
Type=simple
User=onetime
Group=onetime
WorkingDirectory=/var/lib/onetimesecret

# Source environment and start Puma
ExecStart=/bin/bash -lc 'source .env.sh && exec bundle exec puma -C etc/puma.rb'

# Graceful shutdown (systemd sends SIGTERM by default)
TimeoutStopSec=30

# Restart on failure
Restart=on-failure
RestartSec=5

EnvironmentFile=/etc/default/onetimesecret

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/onetimesecret/tmp /var/lib/onetimesecret/log
PrivateTmp=true

[Install]
WantedBy=multi-user.target
