# etc/examples/systemd/onetimesecret-scheduler.service
# Onetime Secret - Scheduler Service
#
# Required only for full authentication mode (PostgreSQL + RabbitMQ).
# Enqueues periodic maintenance jobs (secret cleanup, session expiry, etc.)
#
# Install:
#   sudo cp etc/examples/systemd/onetimesecret-*.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now onetimesecret-scheduler

[Unit]
Description=Onetime Secret Job Scheduler
Documentation=https://docs.onetimesecret.com/en/self-hosting/
After=network.target redis-server.service valkey.service onetimesecret-web.service

[Service]
Type=simple
User=onetime
Group=onetime
WorkingDirectory=/var/lib/onetimesecret

ExecStart=/bin/bash -lc 'source .env.sh && exec bin/ots scheduler'

TimeoutStopSec=30

Restart=on-failure
RestartSec=5

EnvironmentFile=/etc/default/onetimesecret

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/onetimesecret/tmp /var/lib/onetimesecret/log
PrivateTmp=true

[Install]
WantedBy=multi-user.target
