# scripts/s6-rc.d/scheduler/run
#!/command/execlineb -P

# Change to app directory
cd /app

# Preserve environment variables
with-contenv

foreground { s6-echo "[scheduler] Starting job scheduler..." }

# Set RUBY_UNBUFFERED for immediate output
export RUBY_UNBUFFERED 1

exec bin/ots scheduler
