# docker/s6/services/scheduler/finish
#!/command/execlineb -S0

# Log exit for debugging
foreground { s6-echo "[scheduler finish] Scheduler exited: code=${1} signal=${2}" }

# Exit with same code
/bin/sh -c "exit $1" _ ${1}
